# systemd/tunnel-strongswan.service
[Unit]
Description=Tunnel VPN Gateway - strongSwan IPsec (L2TP)
After=network.target network-online.target tunnel-dashboard.service
Wants=network-online.target
PartOf=tunnel.target

[Service]
Type=forking
User=root
ExecStartPre=/usr/bin/node /opt/tunnel-server/src/services/gen-l2tp-config.js
ExecStartPre=/bin/sh -c 'ipsec stop 2>/dev/null || true'
ExecStart=/usr/sbin/ipsec start
ExecStop=/usr/sbin/ipsec stop
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
