# systemd/tunnel-xl2tpd.service
[Unit]
Description=Tunnel VPN Gateway - xl2tpd (L2TP)
After=network.target tunnel-strongswan.service tunnel-dashboard.service
Wants=tunnel-strongswan.service
PartOf=tunnel.target

[Service]
Type=forking
User=root
PIDFile=/run/xl2tpd.pid
ExecStartPre=/usr/bin/node /opt/tunnel-server/src/services/gen-l2tp-config.js
ExecStart=/usr/sbin/xl2tpd -c /opt/tunnel-server/configs/l2tp/xl2tpd.conf -p /run/xl2tpd.pid
ExecStop=/bin/kill -TERM $MAINPID
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
