# systemd/tunnel.target
[Unit]
Description=WireGuard Tunnel Server - All Services (Dashboard + WireGuard)
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/bin/true
ExecReload=/bin/true

[Install]
WantedBy=multi-user.target
