# systemd/tunnel-wireguard.service
[Unit]
Description=WireGuard Tunnel Server - WireGuard Interface
After=network.target network-online.target
Wants=network-online.target
PartOf=tunnel.target

[Service]
Type=oneshot
RemainAfterExit=yes
User=root
ExecStartPre=/bin/sh -c 'wg-quick down wg0 2>/dev/null || true'
ExecStartPre=/usr/bin/node /opt/tunnel-server/src/services/gen-wireguard-config.js
ExecStart=/usr/bin/wg-quick up wg0
ExecStop=/usr/bin/wg-quick down wg0

[Install]
WantedBy=multi-user.target
